# systemd/visitor-counter-mobile-api-user.service
[Unit]
Description=Visitor Counter Mobile Status API
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/raspibob/PersonenZ-hler
Environment=PYTHONPATH=/home/raspibob/PersonenZ-hler/src
EnvironmentFile=%h/.config/personenzaehler/api.env
ExecStart=/home/raspibob/PersonenZ-hler/.venv/bin/python scripts/status_api.py --project-root /home/raspibob/PersonenZ-hler
Restart=on-failure
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
RestrictSUIDSGID=true

[Install]
WantedBy=default.target
