# systemd/visitor-counter-status-api.service
[Unit]
Description=Visitor Counter Local Status API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=raspibob
WorkingDirectory=/home/raspibob/PersonenZ-hler
Environment=PYTHONPATH=/home/raspibob/PersonenZ-hler/src
ExecStart=/home/raspibob/PersonenZ-hler/.venv/bin/python scripts/status_api.py --project-root /home/raspibob/PersonenZ-hler --host 0.0.0.0 --port 8766
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
